# deploy/redis.service
[Unit]
Description=Redis data store for agent-link
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/redis-server /etc/agent-link/redis.conf
ExecStop=/usr/bin/redis-cli shutdown
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
